# etc/crontabs/root
# /etc/crontabs/root:jeff
#
# 1. https://openwrt.org/docs/guide-user/base-system/cron
# 2. https://crontab.guru
#
# FIXME(jeff): I do not believe that this crond daemon has
# support compiled in for sending emails -- cronic does
# indeed work proper, but it is the job of the crond daemon
# to do the actual legwork here. Further research is 
# warranted [1]...
#
# 1. https://openwrt.org/docs/guide-user/services/email/smtp.client#installation3
#
# CONFIRMED(jeff): We must build a Busybox package that has
# sendmail enabled in order for the Busybox implementation
# of cron to send mail for us. While we are at it, if we 
# wind up doing so, we ought to also compile in the 
# run-parts applet. The BASH version we are using is not 
# 100% compatible with ash.
#

PATH=/root/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
SHELL=/bin/ash
MAILFROM=root
MAILTO=i8degrees+ns2@gmail.com
CRON_TZ=America/Chicago

# FIXME(jeff): Eek! run-parts broke my cron jobs. I will put either porting
# or re-writing a replacement for run-parts in the coming weeks...

*/5 * * * * cd / && /etc/periodic/5min/pihole_sync.sh
*/10 * * * * cd / && /etc/periodic/10min
*/15 * * * * cd / && /etc/periodic/15min
*/30 * * * * cd / && /etc/periodic/30min
*/45 * * * * cd / && /etc/periodic/45min

# /etc/cron.hourly is a symlink
0 * * * * cd / && /etc/periodic/hourly

# /etc/cron.daily is a symlink
0 2 * * * cd / && /etc/periodic/daily

# /etc/cron.weekly is a symlink
0 3 * * 6 cd / && /etc/periodic/weekly

# /etc/cron.monthly is a symlink
0 5 1 * * cd / && /etc/periodic/monthly

