# etc/rc.local
#!/bin/sh
#

/etc/init.d/pihole restart

exit 0
